# nanopi/zigbee/zigbee2mqtt.service
[Service]
WorkingDirectory=/home/pi/git/zigbee2mqtt
ExecStart=/usr/bin/node index.js
Restart=always
Environment=NODE_ENV=production
User=pi
